# web/deploy/notion-invoice.service
# /etc/systemd/system/notion-invoice.service
# 配置後: sudo systemctl daemon-reload && sudo systemctl enable --now notion-invoice
[Unit]
Description=BKTSK Notion Invoice (Waku)
After=network.target

[Service]
Type=simple
User=notion-invoice
WorkingDirectory=/opt/notion-invoice/web
# mise のシムを PATH に含め、mise でピン留めした node を使う
Environment=PATH=/home/notion-invoice/.local/share/mise/shims:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
# 秘密は EnvironmentFile で渡す（.env をサーバーのリポジトリに置かない）
EnvironmentFile=/etc/notion-invoice.env
ExecStart=/opt/notion-invoice/web/node_modules/.bin/waku start --host 0.0.0.0 --port 3080
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
